# files/fcommandd.conf
description  "FCommandD Input Server"

start on (filesystem and runlevel [!06])
stop  on runlevel [016]

script
  cd /home/user
  /usr/bin/fcommandd >> /tmp/fcommandd.log 2>&1 
end script

